# openrc/rclocal
#!/sbin/openrc-run
#
# https://github.com/mitchweaver/rclocal
# =======================================

description="bring back /etc/rc.local"

depend()
{
	before dev modules
}

start()
{
	/bin/sh /etc/rc.local
}

